# examples/crontab
# crontab
# 每天早上 9 点执行
41 9 * * * /worker/runner -l /logs -c /worker/.env >/proc/1/fd/1 2>/proc/1/fd/2
